# images/centos7/sshd/init/supervisord.ini
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:nscd]
command=/usr/bin/nscd
